# contrib/velbusd.conf
description	"Velbus daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

script
	# defaults
	SERIAL_PORT="/dev/ttyS0"
	LOGFILE="/var/log/domotica/domotica.log"
	BIND="[::]:8445"

	if [ -r /etc/default/velbusd ]; then
		. /etc/default/velbusd
	fi

	/usr/local/bin/velbusd -f -s "$SERIAL_PORT" -b "$BIND" 2>&1 | /usr/local/bin/velbus-parse.pl >> "$LOGFILE"
end script
