# contrib/velbus.conf
description	"Velbus daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

expect fork
exec /usr/local/bin/velbusd -l /var/log/domotica/domotica.log -b [::]:8445
